# deploy/image/supervisord/supervisord-hadoop01.conf
; supervisord 自含配置 — hadoop01 (全 daemon: ZK+JN+NN+ZKFC+DN+RM+NM+JHS)
; JMX Exporter ports: ZK=10109, JN=10103, NN=10101, DN=10102, RM=10104, NM=10105, JHS=10106, HMS=10110, HS2=10111, HM=10107, RS=10108
[supervisord]
nodaemon=true
user=root
logfile=/logs/supervisord.log
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:zookeeper]
command=/opt/zookeeper/bin/zkServer.sh start-foreground
environment=SERVER_JVMFLAGS="-javaagent:/opt/jmx-exporter/jmx_prometheus_javaagent-0.20.0.jar=10109:/opt/jmx-exporter/config.yml"
autostart=true ; autorestart=false ; startsecs=3  ; priority=10
stdout_logfile=/logs/zk.log  ; stderr_logfile=/logs/zk.err

[program:journalnode]
command=/opt/hadoop/bin/hdfs journalnode
environment=HADOOP_JOURNALNODE_OPTS="-javaagent:/opt/jmx-exporter/jmx_prometheus_javaagent-0.20.0.jar=10103:/opt/jmx-exporter/config.yml"
autostart=true ; autorestart=false ; startsecs=3  ; priority=20
stdout_logfile=/logs/jn.log  ; stderr_logfile=/logs/jn.err

[program:nodemanager]
command=/opt/hadoop/bin/yarn nodemanager
environment=YARN_NODEMANAGER_OPTS="-javaagent:/opt/jmx-exporter/jmx_prometheus_javaagent-0.20.0.jar=10105:/opt/jmx-exporter/config.yml"
autostart=true ; autorestart=false ; startsecs=3  ; priority=40
stdout_logfile=/logs/nm.log  ; stderr_logfile=/logs/nm.err

[program:datanode]
command=/opt/hadoop/bin/hdfs datanode
environment=HADOOP_DATANODE_OPTS="-javaagent:/opt/jmx-exporter/jmx_prometheus_javaagent-0.20.0.jar=10102:/opt/jmx-exporter/config.yml"
autostart=true ; autorestart=false ; startsecs=3  ; priority=30
stdout_logfile=/logs/dn.log  ; stderr_logfile=/logs/dn.err

[program:namenode]
command=/opt/hadoop/bin/hdfs namenode
environment=HADOOP_NAMENODE_OPTS="-javaagent:/opt/jmx-exporter/jmx_prometheus_javaagent-0.20.0.jar=10101:/opt/jmx-exporter/config.yml"
autostart=true ; autorestart=false ; startsecs=10 ; priority=50
stdout_logfile=/logs/nn.log  ; stderr_logfile=/logs/nn.err

[program:zkfc]
command=/opt/hadoop/bin/hdfs zkfc
autostart=true ; autorestart=false ; startsecs=5  ; priority=60
stdout_logfile=/logs/zkfc.log ; stderr_logfile=/logs/zkfc.err

[program:resourcemanager]
command=/opt/hadoop/bin/yarn resourcemanager
environment=YARN_RESOURCEMANAGER_OPTS="-javaagent:/opt/jmx-exporter/jmx_prometheus_javaagent-0.20.0.jar=10104:/opt/jmx-exporter/config.yml"
autostart=true ; autorestart=false ; startsecs=10 ; priority=70
stdout_logfile=/logs/rm.log  ; stderr_logfile=/logs/rm.err

[program:historyserver]
command=/opt/hadoop/bin/mapred historyserver
environment=HADOOP_JOB_HISTORYSERVER_OPTS="-javaagent:/opt/jmx-exporter/jmx_prometheus_javaagent-0.20.0.jar=10106:/opt/jmx-exporter/config.yml"
autostart=true ; autorestart=false ; startsecs=3  ; priority=80
stdout_logfile=/logs/jhs.log ; stderr_logfile=/logs/jhs.err

; ---- Hive (Step 2) ----
[program:hivemetastore]
command=/opt/hive/bin/hive --service metastore
environment=HADOOP_HEAPSIZE="256",JAVA_HOME="/usr/lib/jvm/java-21-openjdk-amd64",HADOOP_OPTS="-javaagent:/opt/jmx-exporter/jmx_prometheus_javaagent-0.20.0.jar=10110:/opt/jmx-exporter/config.yml"
autostart=true ; autorestart=false ; startsecs=10 ; priority=90
stdout_logfile=/logs/hms.log  ; stderr_logfile=/logs/hms.err

[program:hiveserver2]
command=/opt/hive/bin/hiveserver2
environment=HADOOP_HEAPSIZE="256",JAVA_HOME="/usr/lib/jvm/java-21-openjdk-amd64",HADOOP_OPTS="-javaagent:/opt/jmx-exporter/jmx_prometheus_javaagent-0.20.0.jar=10111:/opt/jmx-exporter/config.yml"
autostart=true ; autorestart=false ; startsecs=10 ; priority=100
stdout_logfile=/logs/hs2.log  ; stderr_logfile=/logs/hs2.err

; ---- HBase (Step 2) ----
[program:hmaster]
command=/opt/hbase/bin/hbase master start
environment=HBASE_MASTER_OPTS="-javaagent:/opt/jmx-exporter/jmx_prometheus_javaagent-0.20.0.jar=10107:/opt/jmx-exporter/config.yml",HBASE_HEAPSIZE="256",JAVA_HOME="/usr/lib/jvm/java-8-openjdk-amd64"
autostart=true ; autorestart=false ; startsecs=10 ; priority=110
stdout_logfile=/logs/hm.log  ; stderr_logfile=/logs/hm.err

[program:regionserver]
command=/opt/hbase/bin/hbase regionserver start
environment=HBASE_REGIONSERVER_OPTS="-javaagent:/opt/jmx-exporter/jmx_prometheus_javaagent-0.20.0.jar=10108:/opt/jmx-exporter/config.yml",HBASE_HEAPSIZE="256",JAVA_HOME="/usr/lib/jvm/java-8-openjdk-amd64"
autostart=true ; autorestart=false ; startsecs=10 ; priority=120
stdout_logfile=/logs/rs.log  ; stderr_logfile=/logs/rs.err
